# deploy/qudata-agent.service
[Unit]
Description=QuData Agent Service
After=network.target docker.service auditd.service apparmor.service
Requires=docker.service auditd.service apparmor.service

[Service]
Type=notify

ExecStart=/usr/local/bin/qudata-agent

User=root
Group=root
WorkingDirectory=/opt/qudata-agent

Restart=always
RestartSec=5s

WatchdogSec=30

Environment="QUDATA_API_KEY=YOUR_API_KEY_PLACEHOLDER"

StandardOutput=journal
StandardError=journal
SyslogIdentifier=qudata-agent

[Install]
WantedBy=multi-user.target